# scripts/rc.local
#!/bin/bash
echo "systemd-generators rc.local ok" >> /home/os/p.log
