# scripts/user_timer.service
[Unit]
Decription = Systemd user service for timer
[Service]
Type = oneshot
ExecStart = /home/os/user_timer.sh
[Install]
WantedBy = default.target
